# dating-bot.service
[Unit]
Description=Dating Telegram Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/tgbotSergeyon777
Environment="PATH=/root/tgbotSergeyon777/venv/bin:/usr/bin:/usr/local/bin"
ExecStart=/root/tgbotSergeyon777/venv/bin/python3 /root/tgbotSergeyon777/bot.py
Restart=always
RestartSec=10
StandardOutput=append:/var/log/dating-bot.log
StandardError=append:/var/log/dating-bot-error.log

[Install]
WantedBy=multi-user.target

